# templates/systemd/claude-trigger.timer
[Unit]
Description=Condition check for Claude agent (PROFILE_NAME)
Requires=claude-trigger-PROFILE_NAME.service

[Timer]
OnBootSec=20min
OnUnitActiveSec=30min
RandomizedDelaySec=600
AccuracySec=30s

[Install]
WantedBy=timers.target
