# templates/systemd/claude-trigger.service
[Unit]
Description=Condition check for Claude agent (PROFILE_NAME)

[Service]
Type=oneshot
ExecStart=%h/.local/bin/claude-trigger.sh PROFILE_NAME
TimeoutStartSec=1200
StandardOutput=journal
StandardError=journal
